# run.crontab
*/1 * * * * /bin/sh /git/iaai-checker/ping-okinet.sh
*/1 * * * * /bin/sh /git/iaai-checker/check-lot.sh