# run.crontab
*/5 * * * * /bin/sh /git/iaai-checker/ping-okinet.sh